# kali-config/variant-bsum/includes.chroot/etc/profile.d/no-dotnet-telemetry.sh
export DOTNET_CLI_TELEMETRY_OPTOUT=true
